# .docker/supervisor.conf
[supervisord]
nodaemon=true

[program:schedule-worker]
command=php /var/www/news-backend/artisan schedule:work
autostart=true
autorestart=true
user=www-data
